# services/runit/sv/pipewire-pulse/run
#!/bin/sh

# PipeWire Pulse Audio compatibility runit service script
# Part of Voidance Linux system services

# Source function library if it exists
if [ -r /etc/rc.conf ]; then
    . /etc/rc.conf
fi

# Check if pipewire-pulse binary exists
if ! command -v pipewire-pulse >/dev/null 2>&1; then
    echo "PipeWire Pulse binary not found"
    exit 1
fi

# Create necessary directories
mkdir -p /var/run/pipewire
mkdir -p /etc/pipewire/pipewire-pulse.conf.d

# Set up environment
export PATH="/usr/bin:/usr/sbin:/bin:/sbin"
export XDG_RUNTIME_DIR="${XDG_RUNTIME_DIR:-/run/user/$(id -u)}"

# Start PipeWire Pulse Audio compatibility daemon
echo "Starting PipeWire Pulse Audio compatibility..."
exec pipewire-pulse 2>&1